# scripts/home-athan.service
[Unit]
Description=Home Athan Automation Service
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/Home-Athan
ExecStart=/home/pi/Home-Athan/venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
